# scripts/cronjobs/cron-git-push-notes.cron
5 8-20 * * *  cd ~/notes/ && git add -A && git commit -m "Auto-updating notes" && git push
