# iodined-server-setup/iodined.cron
*/5 * * * * /etc/iodined/keep-iodined-running.sh
